# DEBIAN/postinst
#!/bin/bash
# Rootless (Dopamine hideroot) 安装脚本 — 动态检测 jbroot
set -e

# === 1. 动态检测越狱根路径 ===
detect_jbroot() {
    # a) 环境变量
    if [ -n "$JBROOT" ] && [ -d "$JBROOT" ]; then
        echo "$JBROOT"
        return
    fi

    # b) /var/jb 符号链接 → 解析真实路径
    if [ -L /var/jb ]; then
        local target
        target=$(readlink /var/jb 2>/dev/null)
        if [ -n "$target" ] && [ -d "$target" ]; then
            echo "$target"
            return
        fi
    fi

    # c) /var/jb 目录存在
    if [ -d /var/jb ]; then
        echo "/var/jb"
        return
    fi

    # d) Dopamine hideroot: 扫描 .jbroot-* 目录
    local d
    for d in /private/var/containers/Bundle/Application/.jbroot-*; do
        if [ -d "$d" ]; then
            echo "$d"
            return
        fi
    done

    # e) 回退
    echo "/var/jb"
}

JBROOT=$(detect_jbroot)
echo "[postinst] 检测到 jbroot: $JBROOT"

# === 2. 路径定义 ===
DAEMON_BIN="$JBROOT/usr/bin/ios-autogo"
PLIST="$JBROOT/Library/LaunchDaemons/com.autogo.daemon.plist"
DATA_DIR="/var/mobile/Documents/autogo"

# === 3. Daemon 权限 ===
chmod 755 "$DAEMON_BIN" 2>/dev/null || true
chown root:wheel "$DAEMON_BIN" 2>/dev/null || true

# === 4. 重写 LaunchDaemon plist (修复 jbroot 路径) ===
if [ -f "$PLIST" ]; then
    # 使用 PlistBuddy 或 sed 替换路径
    if command -v PlistBuddy >/dev/null 2>&1; then
        PlistBuddy -c "Set :Program $DAEMON_BIN" "$PLIST" 2>/dev/null || true
        PlistBuddy -c "Delete :ProgramArguments" "$PLIST" 2>/dev/null || true
        PlistBuddy -c "Add :ProgramArguments array" "$PLIST" 2>/dev/null || true
        PlistBuddy -c "Add :ProgramArguments:0 string $DAEMON_BIN" "$PLIST" 2>/dev/null || true
        PlistBuddy -c "Add :ProgramArguments:1 string 8888" "$PLIST" 2>/dev/null || true
    else
        # fallback: sed
        sed -i '' "s|/var/jb/usr/bin/ios-autogo|$DAEMON_BIN|g" "$PLIST" 2>/dev/null || \
        sed -i "s|/var/jb/usr/bin/ios-autogo|$DAEMON_BIN|g" "$PLIST" 2>/dev/null || true
    fi
    echo "[postinst] LaunchDaemon plist 已更新: $DAEMON_BIN"
fi

# === 5. 加载守护进程 ===
/bin/launchctl load "$PLIST" 2>/dev/null || \
"$JBROOT/bin/launchctl" load "$PLIST" 2>/dev/null || true
echo "[postinst] 守护进程已加载"

# === 6. 创建用户数据目录 ===
mkdir -p "$DATA_DIR"/{screenshots,logs,crash,data}
chown mobile:mobile "$DATA_DIR" 2>/dev/null || true

# === 7. 注册 App 图标到 SpringBoard ===
APP_PATH="$JBROOT/Applications/AutoGo.app"
if [ -d "$APP_PATH" ]; then
    chmod -R 755 "$APP_PATH"
    chown -R root:wheel "$APP_PATH" 2>/dev/null || true

    # 重新签名 App
    APP_BIN="$APP_PATH/AutoGo"
    ENT_FILE="$APP_PATH/entitlements.plist"

    # 重写 entitlements.plist 中的 /var/jb/ 路径
    if [ -f "$ENT_FILE" ]; then
        if command -v PlistBuddy >/dev/null 2>&1; then
            PlistBuddy -c "Delete :com.apple.security.exception.files.absolute-path.read-only:0" "$ENT_FILE" 2>/dev/null || true
            PlistBuddy -c "Add :com.apple.security.exception.files.absolute-path.read-only:0 string $JBROOT/" "$ENT_FILE" 2>/dev/null || true
        else
            sed -i '' "s|/var/jb/|$JBROOT/|g" "$ENT_FILE" 2>/dev/null || \
            sed -i "s|/var/jb/|$JBROOT/|g" "$ENT_FILE" 2>/dev/null || true
        fi
        echo "[postinst] entitlements.plist jbroot 已更新"
    fi

    if [ -x "$APP_BIN" ]; then
        if [ -f "$JBROOT/usr/bin/ldid" ]; then
            if [ -f "$ENT_FILE" ]; then
                "$JBROOT/usr/bin/ldid" -S"$ENT_FILE" "$APP_BIN"
                echo "[postinst] App 已使用 entitlements 签名"
            else
                "$JBROOT/usr/bin/ldid" -S "$APP_BIN"
                echo "[postinst] App 已签名 (无 entitlements)"
            fi
        elif command -v ldid >/dev/null 2>&1; then
            if [ -f "$ENT_FILE" ]; then
                ldid -S"$ENT_FILE" "$APP_BIN"
            else
                ldid -S "$APP_BIN"
            fi
            echo "[postinst] App 已使用 ldid 签名"
        else
            echo "[postinst] 警告: ldid 未找到, App 可能无法启动"
        fi
    fi
fi

# === 8. 刷新 uicache ===
"$JBROOT/usr/bin/uicache" -a 2>/dev/null || /usr/bin/uicache -a 2>/dev/null || true

echo "AutoGo Daemon (Rootless) 安装完成!"
echo "  jbroot: $JBROOT"
echo "  Web 控制台: http://设备IP:8888"
echo "  本机访问: http://127.0.0.1:8888"
exit 0
